# recipes-core/systemd/systemd/systemd-boot-update.service.d/override.conf
[Unit]
Requires=efi.mount
After=efi.mount
BindsTo=efi.mount

[Service]
TimeoutStartSec=30s
